# systemd/agent-coder-runner.service
[Unit]
Description=Agent Coder Remote Runner
Wants=network-online.target agent-coder-gateway.service
After=network-online.target agent-coder-gateway.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/Agent-IA-Coder/agent-coder-remote-runner
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
