# systemd/agent-coder-disk-monitor.service
[Unit]
Description=Agent Coder Disk Space Monitor
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/pi/Agent-IA-Coder/agent-coder-disk-monitor
ExecStart=/home/pi/Agent-IA-Coder/agent-coder-disk-monitor/disk_monitor.sh
StandardOutput=journal
StandardError=journal
